# dist/upstart/guid.conf
description     "a small web service that will calculate GUIDs from passed in etkeys"

start on runlevel [2345]
stop on runlevel [!2345]

env PORT="8080"
env HOST="127.0.0.1"
env DAEMON="/opt/guid"

console log

respawn
respawn limit 10 5

exec $DAEMON
